# rtos/modules/core/strip.list
usr/doc/*
usr/share/doc/*
usr/share/locale/*
usr/share/gtk-doc/*
usr/man/*
var/log/packages/*
var/log/removed_packages/*
var/log/removed_scripts/*
var/log/scripts/*
var/log/setup/*
var/lib/dpkg/*
usr/share/vim/vim74/tutor/*
usr/share/vim/vim74/spell/*
usr/share/vim/vim74/lang/*
usr/share/vim/vim74/doc/*
usr/share/i18n/locales/*
usr/share/info/*
etc/rc.d/rc.rpc
etc/rc.d/rc.ntpd
etc/rc.d/rc.messagebus
etc/rc.d/rc.font
etc/rc.d/rc.atop
sbin/fsck.ext2
sbin/fsck.ext3
sbin/fsck.ext4
var/spool/mail/*
